# deploy/com.claudedaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claudedaemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_PATH__</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__HOME__/.claude-daemon/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.claude-daemon/daemon.err</string>
</dict>
</plist>
